# subtitles/root/etc/services.d/jobber/run
#!/usr/bin/with-contenv sh
/jobber/jobberrunner /home/root/.jobber
